# init/mcp-midi-controller.service
[Unit]
Description=mcp-midi-controller — MCP MIDI/OSC rig controller daemon
Documentation=https://github.com/teemow/mcp-midi-controller
# The daemon binds a LAN port for the iPad receiver (auv3_receiver_addr, default
# :7800) so the auv3-probe / AUM app can reach it, and OSC needs the network up.
# Ordering a user unit against system targets is best-effort, so Restart= below
# covers the case where the network is not ready yet when this first starts.
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
# Install the binary first:
#   go install github.com/teemow/mcp-midi-controller/cmd/mcp-midi-controller@latest
# (lands in $GOBIN, or $GOPATH/bin = ~/.go/bin here).
ExecStart=%h/.go/bin/mcp-midi-controller
Restart=on-failure
RestartSec=2
# Two listeners (config.yaml): the MCP streamable-HTTP endpoint on loopback
# (listen_addr, default 127.0.0.1:7799) and the LAN iPad receiver
# (auv3_receiver_addr, default :7800) that ingests auv3-probe dumps, ferries AUM
# session files, and accepts the audio-tap stream. Rig-as-code (devices.yaml,
# device-types/, scenes/, auv3-default-states/) is read from the XDG config dir,
# so no extra environment is needed; stdout/stderr go to the journal.

[Install]
WantedBy=default.target
